# crontab
0 */6 * * * /app/.venv/bin/python -m mitty --verbose
